# deploy/ovh/karaok-backup.timer
[Unit]
Description=Run the KaraOK backup every day

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15m

[Install]
WantedBy=timers.target
