# deploy/ovh/karaok-backup.service
[Unit]
Description=Back up KaraOK MySQL data and audio uploads
After=mysql.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/karaok-backup
PrivateTmp=true
NoNewPrivileges=true
